# overlay/etc/services.d/aspnet/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Run /usr/local/bin/aspnet-start executable.
#======================================================================================================================

aspnet-start
