# overlay/etc/services.d/aspnet/run
#!/usr/bin/with-contenv sh

set -euo pipefail
export BCG_E="aspnet/$(basename ${0})"


#======================================================================================================================
# Run /usr/local/bin/aspnet-start executable.
#======================================================================================================================

aspnet-start
