# overlay/etc/services.d/dotnet/run
#!/usr/bin/with-contenv sh

if [ -z "${DOTNET_RUN_ASSEMBLY}" ] ; then
    echo " - DOTNET_RUN_ASSEMBLY not set."
    exit 1
fi

echo " - starting .NET"
s6-setuidgid www /usr/bin/dotnet_run
