# overlay/etc/services.d/dotnet/run
#!/usr/bin/with-contenv sh

echo " - running ASP.NET application"
s6-setuidgid www dotnet-start-app
