# overlay/etc/services.d/dotnet/run
#!/usr/bin/with-contenv sh
_echo "Starting ASP.NET application"
s6-setuidgid www dotnet-run
