# overlay/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv sh

set -euo pipefail


#======================================================================================================================
# Check for FTP password first
#======================================================================================================================

if [ -z "${FTP_PASSWORD}" ] ; then
    # de-register the vsftpd service
    s6-svc -d /var/run/s6/services/vsftpd
    exit 0
fi


#======================================================================================================================
# Run service
#======================================================================================================================

_echo "Starting vsftp service"
vsftpd
