# overlay/etc/s6-overlay/s6-rc.d/aspnet/finish
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x cron

# Show helpful log message - DO NOT TERMINATE CONTAINER though
def main [...args] { bf-s6 svc finish }
