# overlay/etc/s6-overlay/s6-rc.d/aspnet/run
#!/usr/bin/nu

use bf
use bf-s6
bf env load -x aspnet

# Run ASP.NET application
def main [] {
    # verify ASPNET_ASSEMBLY is set
    let assembly = bf env req ASPNET_ASSEMBLY
    let dir_live = bf env req ASPNET_APP_LIVE

    # verify assembly exists
    let assembly_path = $"($dir_live)/($assembly)"
    if (bf fs is_not_file $assembly_path) {
        bf write notok $"($assembly_path) not found."
        go_to_sleep 30sec
        exit 0
    }

    # start application
    bf write "Starting ASP.NET application."
    bf write debug $" .. ($assembly_path)."

    ^s6-setuidgid www dotnet $assembly_path
}

# Sleep the thread for a set about of time
def go_to_sleep [
    for: duration   # Amount of time to sleep for
] {
    bf write $"Sleeping for ($for)."
    sleep $for
}
