# overlay/etc/s6-overlay/s6-rc.d/aspnet/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-aspnet run ; run preflight"


#======================================================================================================================
# Start dotnet application.
#======================================================================================================================

s6-setuidgid www dotnet ${BF_ASPNET_APP_LIVE}/${BF_ASPNET_ASSEMBLY}
